# linux/etc/systemd/system/backup_root.service
[Unit]
Description=full system backup

[Service]
Nice=19
IOSchedulingClass=2
IOSchedulingPriority=7
ExecStart=/path/to/myBackupScript
