# linux/etc/systemd/system/battery_log.service
[Unit]
Description=battery log

[Service]
Nice=19
IOSchedulingClass=2
IOSchedulingPriority=7
ExecStart=/path/to/myBatteryLog/script
